# .docker/config/crontab
* * * * * /usr/bin/php /var/www/beta/artisan schedule:run >> /dev/null 2>&1